# all-code/packetfs-quantum.container
# PacketFS Quantum Compute Environment - Podman Quadlet
# Runs the INSANE quantum shell with unlimited scaling capabilities

[Unit]
Description=PacketFS Quantum Infinite Compute Shell
After=network-online.target
Wants=network-online.target

[Container]
Image=ubuntu:22.04
ContainerName=packetfs-quantum-shell

# Mount the PacketFS quantum environment
Volume=/.pfs2:/opt/packetfs:Z
Volume=/tmp:/tmp:Z

# Unlimited resources (as much as possible)
Memory=0
MemorySwap=0
CPUQuota=0
PidsLimit=0

# Networking for quantum mesh
Network=packetfs-quantum-network
PublishPort=4000-5000:4000-5000

# Environment for quantum acceleration
Environment=PACKETFS_QUANTUM_MODE=true
Environment=PACKETFS_CORES_UNLIMITED=true
Environment=PACKETFS_REALITY_BREAKING=true
Environment=TERM=xterm-256color

# Privileged mode for maximum power
Privileged=true
SecurityOpt=label=disable

# Keep running the quantum shell
Exec=/bin/bash -c "cd /opt/packetfs && python3 quantum_packetfs_infinite_shell.py"

# Auto-restart if quantum wavefunction collapses
Restart=always
RestartSec=5

[Service]
Type=notify
NotifyAccess=all

[Install]
WantedBy=multi-user.target default.target
